# buildroot/board/teslausb/rootfs_overlay/etc/systemd/system/teslausb.service
[Unit]
Description=TeslaUSB Neo Daemon
After=local-fs.target teslausb-provision.service
Wants=teslausb-provision.service
DefaultDependencies=no

[Service]
Type=notify
ExecStart=/usr/bin/teslausb
Restart=always
RestartSec=2
WatchdogSec=30
NotifyAccess=main
ProtectSystem=strict
# The daemon mounts media under /mnt, hot-adds zram for the snapshot COW via
# /sys/class/zram-control, and creates the dm-snapshot under /dev/mapper. Under
# ProtectSystem=strict these must be writable or snapshot/media-sync fail.
ReadWritePaths=/data /mnt /sys/kernel/config /sys/class/leds /proc/sys/vm /sys/block /sys/class/zram-control /dev/mapper
PrivateTmp=true

[Install]
WantedBy=multi-user.target
